# Docker/rootfs/etc/services.d/automator-server/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid abc /usr/bin/python3 -m flask run
