# Docker/rootfs/etc/services.d/automator-server/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid abc python3 conversion-server/conversion-server/__init__.py
